# templates/50-udev.rules.j2
ACTION=="add", SUBSYSTEM=="scsi" , SYSFS{type}=="0|7|14", \
RUN+="/bin/sh -c 'echo {{ iscsi.timeout }} > /sys$DEVPATH/device/timeout'"
